# scripts/autostart/slapshot.desktop
[Desktop Entry]
Type=Application
Name=Slapshot
Exec=/home/orangepi/bin/slapshot.sh

